# rddl-info.crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin:/root/.local/bin
*/2 * * * *   root	(cd /usr/src/app && python rddl_info/crawl.py --influxdb 1) | logger -t rddl-info
